# novnc/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start app service
# ==============================================================================

bashio::log.info "Starting app..."

exec supervisord -c /app/supervisord.conf