# novnc/app/conf.d/nginx.conf
[program:nginx]
command=nginx
autorestart=true