# novnc/app/conf.d/websockify.conf
[program:websockify]
command=websockify -v --web /usr/share/novnc 8080 localhost:5900
autorestart=true