# novnc/app/conf.d/websockify.conf
[program:websockify]
command=websockify -v --web /usr/share/novnc 127.0.0.1:5901 127.0.0.1:5900
autorestart=true