# novnc/rootfs/etc/services.d/xterm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start xterm service
# ==============================================================================
MESSAGE=$(bashio::config 'message')
bashio::log.info "Starting yterm service..."
xterm -e echo "${MESSAGE}"
