# novnc/rootfs/etc/services.d/fluxbox/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start fluxbox service
# ==============================================================================
bashio::log.info "Starting fluxbox service..."
fluxbox
